# bandonfews.conf
[program:bandonfews]
command=/home/pi/gitwork/bandonfews-nodejs/main.js
autostart=true
autorestart=true
stderr_logfile=/var/log/bandonfews.err.log
stdout_logfile=/var/log/bandonfews.out.log
